# xdg_portal/contrib/systemd/xdg-desktop-portal-pikeru.service.in
[Unit]
Description=Portal service (pikeru file picker implementation)
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.pikeru
ExecStart=@libexecdir@/xdg-desktop-portal-pikeru
Restart=on-failure
